# debian/rmtfs.postinst
#!/bin/sh

#DEBHELPER#

rc-update del rmtfs default || true
rc-update add rmtfs sysinit || true

exit 0
